# files/etc/systemd/system/paccache.service.d/override.conf
[Service]
ExecStart=
ExecStart=/bin/sh -c '/usr/bin/paccache -ruk0 && /usr/bin/paccache -rk2'
