# templates/etc/systemd/system/ddupdate.service.d/serviceopt.conf
[Service]
ExecStart=
ExecStart=/usr/bin/ddupdate --service-option zone={{ dyndns.cf.zone }}
