# files/etc/cron.daily/dehydrated
#!/bin/sh

dehydrated -c 2>&1 | ts %FT%T%z >> /var/log/dehydrated.log
