# files/etc/init.d/weechat
#!/sbin/openrc-run

description="Weechat in tmux"

depend() {
	need net
	use dns netmount
}

start() {
	ebegin "Starting $RC_SVCNAME"
	start-stop-daemon --start --user weechat --name $RC_SVCNAME \
		--start /usr/bin/tmux -- -u2L weechat new-session -d \
			-- /usr/bin/weechat -d /var/lib/weechat
	eend $?
}

status() {
	ebegin "Status of $RC_SVCNAME"
	su - weechat -c '/usr/bin/tmux -L weechat list-sessions'
	eend $?
}

stop() {
	ebegin "Stopping $RC_SVCNAME"
	su - weechat -c '/usr/bin/tmux -L weechat kill-session'
	eend $?
}
