# templates/etc/systemd/system/gunicorn.socket
[Unit]
Description=Gunicorn {{ _key }} socket

[Socket]
ListenStream=/run/gunicorn/{{ _key }}.sock

[Install]
WantedBy=sockets.target
