# templates/etc/systemd/system/getty@tty1.service.d/override.conf
[Service]
ExecStart=
ExecStart=-/usr/sbin/agetty -a {{ autologin.user }} --noclear %I $TERM
